# openode/files/supervisor.conf
{% from "openode/map.jinja" import server with context %}

[program:openode_gunicorn]
command={{ server.dir.base }}/bin/python {{ server.dir.base }}/bin/gunicorn_start --bind={{ server.bind.address }}:{{ server.bind.port }} --workers={{ server.workers }} --bind=unix:{{ server.dir.base }}/openode.sock
directory={{ server.dir.base }}/site
user=openode
autostart=true
autorestart=true
redirect_stderr=true

[program:openode_celery]
command={{ server.dir.base }}/bin/python manage.py celeryd -E --loglevel=INFO
directory={{ server.dir.base }}/site
user=openode
numprocs=1
stdout_logfile=/var/log/openode/celeryd.log
stderr_logfile=/var/log/openode/celeryd.log
autostart=true
autorestart=true
startsecs=3

[group:openode]
programs=openode_gunicorn,openode_celery